# foo.timer
[Unit]
Description=foo's timer file

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
